# bbcetc.cron
#10 * * * * env - `cat /etc/environment` python3 /base/bbcetc.py >> /dev/null 2>&1
10 * * * * env - `cat /etc/environment` python3 /base/bbcetc.py > /var/log/bbcetc.log

